# systemd/server.service.in
[Unit]
Description=Start all the server docker containers
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=${SERVER_PATH:-/srv/server}/
ExecStart=/usr/bin/bash -c 'source ${SERVER_PATH:-/srv/server}/server.sh && server_up'
ExecStop=/usr/bin/bash -c 'source ${SERVER_PATH:-/srv/server}/server.sh && server_down'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
