# systemd/server_sshfs_mount_vault.service.in
[Unit]
Description=Mount Hetzner 'vault' directory with sshfs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target

[Service]
Type=exec

ExecStart=/usr/bin/sshfs -f -p ${STORAGE_SSH_PORT:-22} ${STORAGE_SSH_USER_VAULT:-vault}@${STORAGE_SSH_HOST:-localhost}:/ ${STORAGE_SSH_MOUNT_VAULT:-${SERVER_PATH}/storage/mount/sshfs.vault} -o volname="Hetzner Box - 'vault'",default_permissions,rw,allow_other,idmap=user,uid=${STORAGE_SSH_UID:-0},gid=${STORAGE_SSH_GID:-0},ServerAliveInterval=20,reconnect,ServerAliveCountMax=5,_netdev,auto_cache,dir_cache=yes,StrictHostKeyChecking=no,IdentityFile=${SERVER_PATH}/storage/.ssh/id_ed25519

ExecStop=/usr/bin/umount ${STORAGE_SSH_MOUNT_VAULT:-${SERVER_PATH}/storage/mount/sshfs.vault}

Restart=on-failure
RestartSec=30
