# systemd/server_sshfs_mount_other.service.in
[Unit]
Description=Mount Hetzner 'other' directory with sshfs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target

[Service]
Type=exec

ExecStart=/usr/bin/sshfs -f ${STORAGE_SSH_USER_OTHER:-other}@${STORAGE_SSH_HOST:-localhost}:/ ${STORAGE_SSH_MOUNT_OTHER:-${SERVER_PATH}/storage/mount/sshfs_other} -o default_permissions,rw,allow_other,idmap=user,uid=${STORAGE_SSH_UID:-0},gid=${STORAGE_SSH_GID:-0},ServerAliveInterval=20,reconnect,ServerAliveCountMax=5,_netdev,IdentityFile=${SERVER_PATH}/storage/.ssh/id_ed25519,dir_cache=yes

ExecStop=/usr/bin/umount ${STORAGE_SSH_MOUNT_OTHER:-${SERVER_PATH}/storage/mount/sshfs_other}

Restart=on-failure
RestartSec=30
