# systemd/server_systemctl.service.in
[Unit]
Description=Syncthing - Open Source Continuous File Synchronization
Documentation=man:syncthing(1)
After=network.target server_sshfs_mount_vault.service server_sshfs_mount_other.service server_gocryptfs_mount_vault_private.service server_gocryptfs_mount_vault_generic.service
StartLimitIntervalSec=60
StartLimitBurst=4

[Service]
User={STORAGE_UID:-0}
Group=${STORAGE_GID:-0}

ExecStart=/usr/bin/syncthing serve --no-browser --no-restart --logflags=0 --home="${STORAGE_SYNCTHING_CONFIG_PATH:-${SERVER_PATH}/storage/.syncthing}" --no-default-folder

Restart=on-failure
RestartSec=1
SuccessExitStatus=3 4
RestartForceExitStatus=3 4

# Hardening
ProtectSystem=full
PrivateTmp=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
